# config/99-realsense-d4xx-tegra.rules
# Device rules for Intel RealSense MIPI devices on Jetson

SUBSYSTEM=="d4xx-class", RUN+="/bin/bash -c /usr/local/bin/rs-enum.sh -q > /dev/kmsg'"
